# tools/osbuilder/rootfs-builder/datadog-files/etc/systemd/system/kata-agent.service.d/50-apparmor.conf
[Unit]
After=datadog-apparmor.service
Wants=datadog-apparmor.service

[Service]
AppArmorProfile=usr.bin.kata-agent
