# src/agent/systemd-services/guest-components-sealedsecrets.service

# /etc/systemd/system/guest-services-sealedsecrets.service
[Unit]
Description=Guest Components Sealed Secrets Service
Documentation=https://github.com/kata-containers/kata-containers
PartOf=guest-services.target
Requires=initdata-processor.service
After=initdata-processor.service

[Service]
Type=simple
ExecStart=/usr/bin/guest-services sealed-secrets --socket /run/guest-services/sealedsecrets.sock
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
User=root
Group=root

RuntimeDirectory=guest-services
RuntimeDirectoryMode=0755

[Install]
WantedBy=guest-services.target
