# vendor/github.com/kata-containers/agent/kata-agent.service.in
[Unit]
Description=Kata Containers Agent
Documentation=https://github.com/kata-containers/agent
Wants=kata-containers.target

[Service]
# Send agent output to tty to allow capture debug logs
# from a VM serial port
StandardOutput=tty
Type=simple
ExecStart=@bindir@/@kata-agent@
LimitNOFILE=infinity
ExecStop=/bin/sync ; /usr/bin/systemctl --force poweroff
FailureAction=poweroff
