# src/agent/initdata-processor/systemd-services/initdata-processor.service
# Copyright (c) 2026 Ant Group
#
# SPDX-License-Identifier: Apache-2.0
#
[Unit]
Description=Kata InitData Processor
Documentation=https://github.com/kata-containers/kata-containers
DefaultDependencies=no
After=local-fs.target
Before=kata-agent.service
ConditionPathExists=/usr/bin/initdata-processor
ConditionPathIsExecutable=/usr/bin/initdata-processor

[Service]
Type=oneshot
ExecStart=/usr/bin/initdata-processor
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal
Restart=no
TimeoutStartSec=300
TimeoutStopSec=30

NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/run/measured-cfg

[Install]
WantedBy=multi-user.target