# tools/osbuilder/rootfs-builder/ubuntu/superprotocol/check-config-files.service
[Unit]
Description=Check and Copy Config Files Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check-config-files.sh
