# tools/osbuilder/rootfs-builder/ubuntu/superprotocol/local-registry.service
[Unit]
Description=Run local registry & fileserver by using Hauler.

[Service]
Type=simple
ExecStart=/usr/local/bin/local-registry.sh

[Install]
WantedBy=multi-user.target
